# files/etc/supervisord.d/livy-server.conf
[program:livy-server]
command=livy-server
startsecs=2
stopwaitsecs=10
user=livy
redirect_stderr=true
stdout_logfile=/var/log/livy/livy-server.log
autostart=true
stopasgroup=true
killasgroup=true